# deploy/systemd/robot-game-portal-online.timer.example
[Unit]
Description=Run robot-game portal online count sender every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
Unit=robot-game-portal-online.service

[Install]
WantedBy=timers.target

